# docker/config/crontab
0 0 * * * python3 /reportEngine/pj_remates.py
