# docker/config/crontab
0 0 * * * python3 /home/tony/REMATES/DOKER/pj_remates.py >> /var/log/pj_remates.log
